# crontab
#run python script every 15 minutes
*/15 * * * * /usr/local/bin/python /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2